# systemd/portwaechter-v2.service
[Unit]
Description=PortWächter V2 Hedge-Fund Edition
After=network-online.target

[Service]
Type=oneshot
User=ascent
Group=ascent
WorkingDirectory=/opt/portwaechter
EnvironmentFile=/etc/portwaechter/portwaechter.env
ExecStart=/opt/portwaechter/.venv/bin/python -m modules.v2.main run
