# systemd/portwaechter-portfolio.service
[Unit]
Description=PortWächter Portfolio Ingest
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
WorkingDirectory=/opt/portwaechter
ExecStart=/opt/portwaechter/.venv/bin/python -m modules.portfolio_ingest.main run

[Install]
WantedBy=multi-user.target
